# templates/openobserve.service.j2
[Unit]
Description=The OpenObserve server
After=syslog.target network-online.target remote-fs.target nss-lookup.target


[Service]
Type=simple
LimitNOFILE=65535
EnvironmentFile={{ openobserve_dir }}/openobserve.env
ExecStart={{ openobserve_dir }}/openobserve
ExecStop=/bin/kill -s QUIT $MAINPID
Restart=on-failure

DynamicUser=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
StateDirectory=openobserve
CacheDirectory=openobserve
WorkingDirectory={{ openobserve_dir }}
NoNewPrivileges=true
MountFlags=slave

ProtectClock=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectHostname=true
ProtectKernelTunables=true
RestrictRealtime=true
ProtectProc=invisible
ProtectKernelLogs=true
RestrictNamespaces=true
CapabilityBoundingSet=~CAP_SETUID CAP_SETGID CAP_SETPCAP
CapabilityBoundingSet=~CAP_CHOWN CAP_FSETID CAP_SETFCAP
CapabilityBoundingSet=~CAP_SYS_PTRACE
CapabilityBoundingSet=~CAP_NET_ADMIN
CapabilityBoundingSet=~CAP_SYS_ADMIN
PrivateUsers=true
RestrictAddressFamilies=~AF_UNIX
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
